# Conception/main-server/supervisord.conf
# ================================================================
#  SUPERVISOR CONFIG — main-server
#  Manages: Cowrie + Mosquitto + RIoTPot + HTTP Server
#  PFE 2025 — Mohamed Bouzira & Bakhti Rayane Abderaouef
# ================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
user=root
nodaemon=true                         ; run in foreground (required for Docker)
logfile=/dev/null                     ; all logs go to service .json files
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid


# ----------------------------------------------------------------
# COWRIE — SSH/Telnet honeypot
# runs as user cowrie (never root)
# ----------------------------------------------------------------
[program:cowrie]
command=/bin/bash -c "/home/cowrie/cowrie/cowrie-env/bin/twist cowrie"
directory=/home/cowrie/cowrie
user=cowrie
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null

# ----------------------------------------------------------------
# MOSQUITTO — MQTT broker
# ----------------------------------------------------------------
[program:mosquitto]
command=mosquitto -c /etc/mosquitto/mosquitto_custom.conf
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null

# ----------------------------------------------------------------
# CoAP SERVER — custom Python server
# ----------------------------------------------------------------
[program:coap-server]
command=python3 /opt/coap_server.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null
environment=PYTHONUNBUFFERED="1"

# ----------------------------------------------------------------
# IOT PLATFORM — Flask dashboard (port 80)
# ----------------------------------------------------------------
[program:iot_platform]
command=python3 /opt/IoT_platform/app.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null
environment=PYTHONUNBUFFERED="1"

# ----------------------------------------------------------------
# TSHARK — network attack detector (SYN flood, port scan, ICMP)
# Runs as a standalone Python daemon — no shell wrapper needed.
# Its own retry loop handles tshark restarts internally.
# ----------------------------------------------------------------
[program:tshark-capture]
command=python3 /opt/tshark_to_json.py
autostart=true
autorestart=unexpected
startretries=3
startsecs=5
stdout_logfile=/dev/null
stderr_logfile=/dev/null
environment=PYTHONUNBUFFERED="1",CAPTURE_IFACE="eth0"

# ----------------------------------------------------------------
# MQTT LOGGER — subscribes to all topics, logs full payload
# ----------------------------------------------------------------
[program:mqtt-logger]
command=python3 /opt/mqtt_logger.py
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/null
stderr_logfile=/dev/null
environment=PYTHONUNBUFFERED="1"

# ----------------------------------------------------------------
# HTTP SERVER — Flask (receives data from camera + doorbell)
# ----------------------------------------------------------------
[program:http-server]
command=python3 /opt/http_server.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null
environment=PYTHONUNBUFFERED="1"

# ----------------------------------------------------------------
# LOG ROTATION — runs every hour to rotate large JSON log files
# uses copytruncate: copies then truncates open files safely
# no null bytes, no broken pipes
# ----------------------------------------------------------------
[program:logrotate]
command=/bin/bash -c "while true; do logrotate /etc/logrotate.d/honeypot; sleep 3600; done"
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null